# packaged/etc/supervisord.d/statsd.ini
[program:statsd]
command=/usr/local/bin/start_statsd.sh
redirect_stderr=true
#stdout_logfile = /var/log/supervisor/%(program_name)s.log
stdout_syslog=true

